# install/supervisor/webserver.conf
[program:webserver]
command=/usr/local/bin/run.sh
redirect_stderr=true
stdout_logfile=/log/webserver_server.log
